# kindle/upstart/kindle-dashboard.conf
description "Open Kindle Planner native dashboard"

start on started lab126_gui
stop on stopping lab126_gui

pre-start script
  test -x /mnt/us/documents/kindle-dashboard-launch.sh || { stop; exit 1; }
end script

exec /mnt/us/documents/kindle-dashboard-launch.sh
